# server_deploy/voiceflow-server.service
[Unit]
Description=Voice Flow License Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/voice-flow-server
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
